# tpm/aziot-tpmd/aziot-tpmd.socket.in
[Unit]
Description=Azure IoT TPM Service API socket
PartOf=aziot-tpmd.service

[Socket]
ListenStream=@socket_dir@/tpmd.sock
SocketMode=0660
DirectoryMode=0755
SocketUser=@user_aziottpm@
SocketGroup=@user_aziottpm@

[Install]
WantedBy=sockets.target
