# cert/aziot-certd/aziot-certd.socket.in
[Unit]
Description=Azure IoT Certificates Service API socket
PartOf=aziot-certd.service

[Socket]
ListenStream=@socket_dir@/certd.sock
SocketMode=0660
DirectoryMode=0755
SocketUser=@user_aziotcs@
SocketGroup=@user_aziotcs@

[Install]
WantedBy=sockets.target
